# .config/systemd/user/mpd.service
[Unit]
Description=Music Player Daemon

[Service]
Type=forking
ExecStart=/usr/bin/mpd

[Install]
WantedBy=default.target
