# .config/systemd/user/plasma-session.target
[Unit]
Description=Plasma session
BindsTo=graphical-session.target
PartOf=sway-session.target x-session.target
