# .config/systemd/user/plasma-session.target
[Unit]
Description=session started directly by Plasma
BindsTo=graphical-session.target
PartOf=sway-session.target xmonad-session.target
