# .config/systemd/user/sway-session.target
[Unit]
Description=session started directly by Sway
BindsTo=sway.service
