# .config/systemd/user/sway.service
[Unit]
Description=Sway Watcher

[Service]
Type=forking
ExecStart=/bin/sh -c '~/bin/wm-wait sway'
ExecReload=/bin/true
