# .config/systemd/user/sway-session.target
[Unit]
Description=Sway session
BindsTo=wayland-session.target
