# .config/systemd/user/sway-session-pre.target
[Unit]
Description=Wayland session startup
RefuseManualStart=yes
StopWhenUnneeded=yes
