# .config/systemd/user/sway.service
[Unit]
Description=Sway Watcher

[Service]
ExecStart=$HOME/bin/wm-wait sway
ExecReload=/bin/true
