# .config/systemd/user/sway-session.target
[Unit]
Description=Sway session
BindsTo=graphical-session.target wayland-session.target
