# .config/systemd/user/sway.service
[Unit]
Description=Sway watcher
RefuseManualStart=yes

[Service]
ExecStart=/usr/bin/swaymsg -q --monitor -t subscribe '["tick"]'
SuccessExitStatus=1
