# .config/systemd/user/sway-session.target
[Unit]
Description=Sway session
PartOf=graphical-session.target
