# .config/systemd/user/taffybar-status.service
[Unit]
Description=Taffybar status widgets
BindsTo=taffybar.service

[Service]
ExecStart=%h/bin/taffybar-status

[Install]
WantedBy=taffybar.service
