# .config/systemd/user/taffybar.service
[Unit]
Description=Taffybar
PartOf=xmonad-session.target

[Service]
ExecStart=/usr/local/bin/taffybar

[Install]
WantedBy=xmonad-session.target
