# .config/systemd/user/urxvtd.service
[Unit]
Description=urxvt terminal daemon
PartOf=x-session.target

[Service]
ExecStart=/usr/bin/urxvtd

[Install]
WantedBy=x-session.target
