# .config/systemd/user/waybar.service
[Unit]
Description=Waybar
PartOf=sway-session.target

[Service]
ExecStart=/usr/bin/waybar

[Install]
WantedBy=sway-session.target
