# .config/systemd/user/waybar-sway-workspaces.service
[Unit]
Description=Sway workspaces widgets for Waybar
PartOf=graphical-session.target

[Service]
ExecCondition=/usr/bin/test -S ${SWAYSOCK}
ExecStart=/bin/sh -c '~/bin/waybar-sway-workspaces 10'

[Install]
WantedBy=graphical-session.target
