# .config/systemd/user/wayland-session.target
[Unit]
Description=Wayland session
RefuseManualStart=yes
StopWhenUnneeded=yes
