# .config/systemd/user/x-resources.service
[Unit]
Description=load X resources
PartOf=graphical-session.target

[Service]
Type=oneshot
ExecStart=/usr/bin/xrdb -merge %h/.Xresources
RemainAfterExit=yes

[Install]
WantedBy=graphical-session.target
