# .config/systemd/user/x-session.target
[Unit]
Description=X session
PartOf=graphical-session.target
