# .config/systemd/user/x-session.target
[Unit]
Description=X session
BindsTo=x11.service
