# .config/systemd/user/x11.service
[Unit]
Description=X11 watcher

[Service]
ExecStart=/usr/bin/xprop -root -spy
StandardOutput=null
