# .config/systemd/user/x11.service
[Unit]
Description=X11 watcher
RefuseManualStart=yes

[Service]
ExecStart=/usr/bin/xprop -root -spy
StandardOutput=null
SuccessExitStatus=1
