# .config/systemd/user/xmonad-session.target
[Unit]
Description=The session started directly by xmonad
BindsTo=xmonad.service
