# .config/systemd/user/xmonad-session.target
[Unit]
Description=session started directly by xmonad
BindsTo=xmonad.service
