# .config/systemd/user/xmonad-session-pre.target
[Unit]
Description=xmonad session startup
RefuseManualStart=yes
StopWhenUnneeded=yes
