# .config/systemd/user/xmonad-session.target
[Unit]
Description=xmonad session
PartOf=plasma-session.target
