# .config/systemd/user/xmonad.service
[Unit]
Description=xmonad Watcher

[Service]
Type=forking
ExecStart=/bin/sh -c '~/bin/wm-wait xmonad'
ExecReload=/bin/true
