# .config/systemd/user/xmonad.service
[Unit]
Description=xmonad Watcher

[Service]
ExecStart=$HOME/bin/wm-wait xmonad
ExecReload=/bin/true
