# .config/autostart/plasma-session.desktop
[Desktop Entry]
Exec=/bin/systemctl --user restart plasma-session.target
Icon=plasma
Name=Plasma Session
Type=Application
X-KDE-AutostartScript=true
