# .config/systemd/user/barrierc@hawk.service.d/10-host-condition.conf
[Unit]
ConditionHost=|falcon
ConditionHost=|owl
